# documentation/supervisor/jibupro_celerybeat.conf
[program:jibuprocelerybeat]
; Set full path to celery program if using virtualenv
command=/home/moshthepitt/.virtualenvs/jibupro/bin/celerybeat -A template --loglevel=INFO

; remove the -A myapp argument if you are not using an app instance

directory=/srv/jibupro
user=moshthepitt
numprocs=1
stdout_logfile=/var/log/celery/answers_beat.log
stderr_logfile=/var/log/celery/answers_beat.log
autostart=true
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
